# .config/systemd/user/dotfiles.timer
[Unit]
Description=dotfiles update timer

[Timer]
Type=oneshot
OnCalendar=*-*-* *:00/10:00
Persistent=true

[Install]
WantedBy=timers.target
